# linux/rootfs/board/common/rootfs_overlay/usr/lib/systemd/system/emiso.service
[Unit]
Description=Starting Emiso-engine

[Service]
Restart=always
ExecStart=/usr/local/bin/emiso_engine

[Install]
WantedBy=multi-user.target
Enabled=true






