# linux/rootfs/board/common/rootfs_overlay/usr/lib/systemd/system/emiso_startup.service
[Unit]
Description=EMISO startup Service

[Service]
Type=oneshot
ExecStart=/etc/emiso/startup.sh

[Install]
WantedBy=multi-user.target
Enabled=true
