# linux/rootfs/board/common/rootfs_overlay/usr/lib/systemd/system/remove_logs.service
[Unit]
Description=Clean SOO Logs

[Service]
Type=oneshot
ExecStart=/etc/soo/remove_logs.sh

[Install]
WantedBy=multi-user.target
Enabled=true
