# linux/rootfs/board/common/rootfs_overlay/usr/lib/systemd/system/soo_startup.service
[Unit]
Description=SOO startup Service

[Service]
Type=oneshot
ExecStart=/etc/soo/startup.sh

[Install]
WantedBy=multi-user.target
Enabled=true
