# agency/rootfs/board/merida/rootfs_overlay/etc/udev/rules.d/00-agency.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="mlan0", RUN+="/root/agency.sh"

