# agency/rootfs/board/rpi4_64/rootfs_overlay/etc/udev/rules.d/00-agency.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan1", RUN+="/sbin/ifconfig wlan1 up"

