# linux/rootfs/board/rpi4_64/rootfs_overlay/usr/lib/systemd/system/syslog-ng@.service
[Unit]
Description=System Logging Service

[Service]
ExecStart=/usr/sbin/syslog-ng -F
StandardOutput=journal
StandardError=journal
Restart=on-failure

[Install]
Alias=syslog.service
WantedBy=multi-user.target
Enabled=true
